# Tools/Drone/etc/systemd/system/aether-1.service
[Unit]
Description=aether script

[Service]
ExecStart=/bin/bash /usr/local/bin/aether-1.sh &

[Install]
WantedBy=multi-user.target